# etc/supervisor.d/weblogic.conf.ini
[program:weblogic]
user=camunda
command=%(ENV_WLS_DOMAIN_HOME)s/startWebLogic.sh
startretries=0
autorestart=false
autostart=false
