# etc/supervisor.d/weblogic.conf.ini
[program:weblogic]
user=camunda
command=/home/camunda/oracle/weblogic/domains/mydomain/startWebLogic.sh
startretries=0
autorestart=false
autostart=false
