# etc/supervisor.d/x11vnc.conf.ini
[program:x11vnc]
command=/bin/bash -c "x11vnc -display ${DISPLAY} -forever"
