# etc/supervisor.d/xvfb.conf.ini
[program:xvfb]
command=/bin/bash -c "Xvfb ${DISPLAY} -shmem -screen 0 ${DISPLAY_WIDTH}x${DISPLAY_HEIGHT}x${DISPLAY_DEPTH}"
